# files/gh-ar.service.j2
[Unit]
Description=GitHub Actions Runner (self-hosted) on Raspberry Pi
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/actions-runner          # update to your runner path
ExecStart=/usr/bin/bash -lc '/home/pi/actions-runner/./run.sh'
Restart=always
RestartSec=60s
KillMode=process
Environment=RUNNER_ALLOW_RUNASROOT=0
# Optional: limit resources
# MemoryMax=300M
# CPUQuota=50%

[Install]
WantedBy=multi-user.target
